# chronos/monkeyscrape.cron

0 8,16 * * * /usr/local/bin/monscrape -d -g > /proc/1/fd/1 2>/proc/1/fd/2


